# launchd/com.shakacode.agent-coord-refresh.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shakacode.agent-coord-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "__AGENT_COORD_REPO__" &amp;&amp; git fetch --quiet origin state &amp;&amp; git checkout --quiet FETCH_HEAD</string>
  </array>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/agent-coord-refresh.out</string>

  <key>StandardErrorPath</key>
  <string>/tmp/agent-coord-refresh.err</string>
</dict>
</plist>
