# launchd/com.shakacode.agent-coord-heartbeat.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shakacode.agent-coord-heartbeat.__AGENT_ID__</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>set -a; . "__AGENT_COORD_ENV_FILE__"; : "${AGENT_COORD_API_URL:?set in __AGENT_COORD_ENV_FILE__}"; : "${AGENT_COORD_API_TOKEN:?set in __AGENT_COORD_ENV_FILE__}"; set +a; cd "__AGENT_COORD_REPO__" &amp;&amp; bin/agent-coord heartbeat --agent-id "__AGENT_ID__" --repo "__TARGET_REPO__" --target "__TARGET__" --batch-id "__BATCH_ID__" --branch "__BRANCH__" --status "${AGENT_COORD_STATUS:-in_progress}"</string>
  </array>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/agent-coord-heartbeat.__AGENT_ID__.out</string>

  <key>StandardErrorPath</key>
  <string>/tmp/agent-coord-heartbeat.__AGENT_ID__.err</string>
</dict>
</plist>
